# scripts/defab_finetune.service
[Unit]
Description=DeFAb Fine-Tuning (Spot-VM Resilient)
Documentation=https://github.com/anonymous-authors/defab/blob/main/hpc/AZURE_FINETUNE.md
After=network-online.target nvidia-fabricmanager.service
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/blanc

# /etc/defab/secrets holds HF_TOKEN and any API keys used by eval steps.
# Create with:
#   sudo mkdir -p /etc/defab && sudo chmod 700 /etc/defab
#   echo 'HF_TOKEN=hf_...' | sudo tee /etc/defab/secrets >/dev/null
#   sudo chmod 600 /etc/defab/secrets
EnvironmentFile=-/etc/defab/secrets
Environment="HF_HOME=/data/hf_cache"
Environment="TRANSFORMERS_CACHE=/data/hf_cache"
Environment="RESULTS_BASE=/data/defab_results"
Environment="REPO_DIR=/home/azureuser/blanc"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONFAULTHANDLER=1"
Environment="NCCL_DEBUG=WARN"
Environment="NCCL_ASYNC_ERROR_HANDLING=1"
Environment="TOKENIZERS_PARALLELISM=false"
# systemd strips the login shell's PATH; put conda env + conda base ahead of
# system paths so python, torchrun, vllm, etc. resolve. The orchestrator also
# sources conda.sh at startup (belt-and-suspenders) but this line means all
# children -- including ones that bypass the orchestrator wrapper -- have a
# sane PATH from the moment systemd starts the unit.
Environment="DEFAB_CONDA_ENV=defab"
Environment="PATH=/home/azureuser/miniconda3/envs/defab/bin:/home/azureuser/miniconda3/bin:/home/azureuser/anaconda3/envs/defab/bin:/home/azureuser/anaconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Favor the orchestrator over torchrun children if the OOM killer is triggered
# -- we want the orchestrator to survive and mark the step as failed rather
# than silently dying mid-write.
OOMScoreAdjust=-500

ExecStart=/usr/local/bin/defab_finetune

# NO ExecStop. systemd already SIGTERMs MainPID itself (because KillSignal=SIGTERM
# below) on `systemctl stop` -- adding `ExecStop=/bin/kill -SIGTERM $MAINPID` is
# both redundant and actively harmful: by the time ExecStop runs, our IMDS-poller-
# triggered SIGTERM handler has already exit 0'd the orchestrator and $MAINPID is
# empty. Bare `kill` then prints its usage and exits 1, causing systemd to
# mis-classify the graceful spot-preemption shutdown as a Failure, which in turn
# fires Restart=on-failure on a VM that Azure is about to evict (wasting restart
# budget from StartLimitBurst).

# Propagate SIGTERM to the whole cgroup so torchrun ranks, vLLM, and data
# prep subprocesses all see it. 45s stop timeout leaves 15s buffer beyond the
# Azure eviction window for final checkpoint flush + state-file fsync.
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
TimeoutStopSec=45

# on-failure restarts after crashes. The orchestrator exits 0 on SIGTERM and
# does NOT retrigger (the VM is being evicted; systemd will start us fresh on
# next boot via WantedBy=multi-user.target).
Restart=on-failure
RestartSec=30

# Ample fd limit for torchrun + many-rank NCCL + HuggingFace Hub + vLLM.
LimitNOFILE=1048576

# Journald captures structured stdout/stderr. Combined with tee in the script
# this gives us redundancy: log survives in both journald and $RESULTS_BASE/run.log.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=defab_finetune

[Install]
WantedBy=multi-user.target
